# docker-rootfs/sim/etc/supervisor/supervisord.conf
[supervisord]
nodaemon=true ; run supervisord as a foreground process
logfile=/tmp/supervisord.log
loglevel=error
pidfile=/tmp/supervisord.pid

[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock         ; use a unix:// URL  for a unix socket

[program:lxpanel]
priority=15
; The sleep is to give time for the display to start
command=bash -c "sleep 5; vglrun /usr/bin/lxpanel --profile LXDE"
autostart=true
autorestart=true
environment=DISPLAY=":1",VGL_DISPLAY=egl
stdout_logfile=/tmp/lxpanel.log
redirect_stderr=true

[program:pcmanfm]
priority=15
; The sleep is to give time for the display to start
command=bash -c "sleep 5; vglrun /usr/bin/pcmanfm --desktop --profile LXDE"
autostart=true
autorestart=true
environment=DISPLAY=":1",VGL_DISPLAY=egl
stdout_logfile=/tmp/pcmanfm.log

[program:turbovnc]
priority=10
command=/turbovnc_wrapper.sh -geometry '1440x1080' -vgl -noreset -SecurityTypes None -rfbport 5900 -fg -log /dev/stdout
autostart=true
autorestart=true
stopsignal=INT
environment=DISPLAY=":1",VGL_DISPLAY=egl,TVNC_WM=openbox-session
stdout_logfile=/tmp/turbovnc.log
redirect_stderr=true

[group:vnc]
# turbovnc configuration
programs=lxpanel,pcmanfm,turbovnc
